# dp3/template/supervisor/supervisord.conf.d/workers.ini
[program:w] ; workers
# add -v to enable debug prints
command = {{DP3_EXE}} worker -v {{DP3_APP}} {{CONFIG_DIR}} %(process_num)d

priority = 10

; Run several workers in parallel
numprocs = {{WORKER_COUNT}}
; WARNING: If changing number of worker processes, the following process 
; must be followed:
; 1. stop all inputs (e.g. warden_receiver, updater)
; 2. when all queues are empty, stop all workers
; 3. reconfigure queues in RabbitMQ using /nerd/scripts/rmq_reconfigure.sh
; 4. change the settings here and in nerd.yml
; 5. reload supervisord and start everything again

process_name = worker%(process_num)d

; ** (Re)starting **
; Wait for 5 seconds until program is considered sucessfully running
startsecs = 5
startretries = 1

; Automatically restart if program exits with an exit code other than 0
autorestart = unexpected
exitcodes = 0

; Give program 30 seconds to exit before it's killed
stopwaitsecs = 30

; ** Logging **
; Redirect stderr to stdout (results in just one log file)
redirect_stderr = true

stdout_logfile_maxbytes = 200MB
stdout_logfile_backups = 5

stdout_logfile = /var/log/{{DP3_APP}}/worker%(process_num)d.log

; prevent stdout to be buffered, otherwise it's printed with a long delay
environment = PYTHONUNBUFFERED=1
